# scripts/systemd/mcgcraft.socket
[Socket]
ListenStream=8080
FlushPending=yes
ListenFIFO=%t/mcgcraft.stdin

[Install]
WantedBy=default.target
