# scripts/systemd/mcgcraft.socket
[Socket]
ListenStream=8080
FlushPending=yes
#ListenFIFO=%t/mcgcraft.stdin
#FlushPending=yes

[Install]
WantedBy=default.target
